# bin/compile
#!/bin/sh
#
# chkconfig: 35 99 99
# description: Node.js /home/nodejs/node-server/server.js
#

if [ $1 ] ; then
	Configs=$1
else
	Configs="`dirname ${0}`/../conf/fexy.conf"
fi

. $Configs

rm -rf "$AppPath/"
cp -r "$SrcPath/" "$AppPath/"
$Coffee -b -c $AppPath
find "$AppPath/" -name '*.coffee' -exec rm {} \;
