# bin/fexy
#!/bin/sh
#
# chkconfig: 35 99 99
# description: Node.js /home/nodejs/node-server/server.js
#

NodeJS="/usr/local/bin/node"
App="`dirname ${0}`/../fexy.js"
lockfile="/var/lock/subsys/node-server"

start()
{
    if [ ! -f $lockfile ] ; then
        echo -n "Starting $App: "
        su -c "nohup $NodeJS $App start >> /dev/null 2>&1 &" && echo "SUCCESS" || echo "FAILED"
        RETVAL=$?
        [ $RETVAL -eq 0 ] && touch $lockfile
    else
        echo "$App is locked."
        RETVAL=1
    fi
}
 
stop()
{
    echo -n "Stopping $App: "
    pid=`ps -aefw | grep "$NodeJS $App" | grep -v " grep " | awk '{print $2}'`
    kill -9 $pid > /dev/null 2>&1 && echo "SUCCESS" || echo "FAIL"
    RETVAL=$?
    [ $RETVAL -eq 0 ] && rm -f $lockfile
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
    	$NodeJS $App $@
        RETVAL=$?
esac

exit $RETVAL
