# bin/fexy
#!/bin/sh
#
# chkconfig: 35 99 99
# description: Node.js /home/nodejs/node-server/server.js
#

Configs="`dirname ${0}`/../conf/fexy.conf"

if [ $2 ] ; then
	Configs=$2
else
	Configs="`dirname ${0}`/../conf/fexy.conf"
fi

. $Configs

. $Configs
 
start()
{
    if [ ! -f $lockfile ] ; then
        echo -n "Starting $App: "
        su -c "umask 0000"
        su -c "nohup $NodeJS $App >> $LogFile 2>&1 &" && echo "SUCCESS" || echo "FAILED"
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch $lockfile
    else
        echo "$App is locked."
        RETVAL=1
    fi
}
 
stop()
{
    echo -n "Stopping $App: "
    pid=`ps -aefw | grep "$NodeJS $App" | grep -v " grep " | awk '{print $2}'`
    kill -9 $pid > /dev/null 2>&1 && echo "SUCCESS" || echo "FAIL"
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f $lockfile
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        RETVAL=1
esac

exit $RETVAL
